# 25_SystemD_Service_Recipe/systemd-example/files/sysd.service
[Unit]
Description=Continuous Hello World Service

[Service]
ExecStart=/bin/bash /usr/bin/sysd.sh

[Install]
WantedBy=multi-user.target